# cups/rootfs/etc/services.d/avahi/run
#!/usr/bin/with-contenv bashio

exec avahi-daemon